# debian/shake-cursor/usr/lib/systemd/user/shake-cursor.service
[Unit]
Description=Shake-to-Find cursor enlarger (user service)
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/shake-cursor
Restart=on-failure

[Install]
WantedBy=default.target
